# .docksal/services/cli/crontab
* * * * * bash -l -c '/usr/local/bin/php /var/www/bin/console drupguard:cron --cron-only'
